# .config/systemd/user/memcached.service
[Unit]
Description=Memcached Server

[Service]
Type=simple
ExecStart=/usr/bin/memcached
ExecStop=/bin/kill -15 $MAINPID
