# ansible/roles/system/files/paccache.service
[Unit]
Description=Pacman Cache Cleaner

[Service]
Type=simple
ExecStart=/usr/bin/paccache -ruk0
