# ansible/roles/system/files/paccache.timer
[Unit]
Description=Clean Pacman Cache Weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
