# ansible/roles/system/files/pacoptimize.timer
[Unit]
Description=Optimize Pacman Database Weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
